# attachments/mbsync.service.timer
[Unit]
Description=mbsync timer

[Timer]
OnCalendar=*:0/3
Unit=mbsync.service

[Install]
WantedBy=default.target
